# config/crontab
# mercstoria
*/5 * * * * /usr/bin/php /home/mercstoria/task/quest.php &>> /home/mercstoria/log 
0,30 * * * * /usr/bin/php /home/mercstoria/task/login_bonus.php &>> /home/mercstoria/log

