# config/crontab
# mercstoria
*/10 * * * * /usr/bin/php /home/mercstoria/task/quest.php 
30 * * * * /usr/bin/php /home/mercstoria/task/login_bonus.php